# src/joycond_amd64/DEBIAN/postinst
#!/bin/bash
systemctl enable --now joycond
